# Optimized-modules/shipping.service
[Unit]
Description=Shipping Service

[Service]
User=roboshop
// highlight-start
Environment=CART_ENDPOINT=cart.devopsmaster.xyz:8080
Environment=DB_HOST=mysql.devopsmaster.xyz
// highlight-end
ExecStart=/bin/java -jar /app/shipping.jar
SyslogIdentifier=shipping

[Install]
WantedBy=multi-user.target
